# ddns-updater.timer
[Unit]
Description=Run ddns-updater.service once daily

[Timer]
OnCalendar=*:00:00