# ddns-updater.service
[Unit]
Description=Updates Dynamic DNS records.

[Service]
Type=oneshot
ExecStart=/usr/bin/local/ddns-updater