# ddns-updater.timer
[Unit]
Description=Run ddns-updater once an hour to limit downtime

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target